# arch-maint.timer
# systemd timer file to run arch-maint service every Saturday at
# 6.45am. If the machine was down during this time, run at next restart.

[Unit]
Description=Runs arch-maint.sh every week

[Timer]
OnCalendar=Sat 06:45
Persistent=true

[Install]
WantedBy=timers.target
